# .github/workflows/scripts/dockerd.service
[Service]
Environment="HTTP_PROXY=http://localhost:9"
Environment="HTTPS_PROXY=https://localhost:9"
Environment="NO_PROXY=localhost,127.0.0.1"
